# bash/task4/tmp/timer.service
[Unit]
Description=write time every second

[Service]
Type=oneshot
StartLimitBurst=0
ExecStart=/home/ilya/RubyLearn/bash/task4/task4.sh

