# bash/task4/tmp/timer.timer
[Unit]

[Timer]
AccuracySec=1us
OnCalendar=*:*:0/1
Unit=timer.service

[Install]
WantedBy=timers.target
